# datadog/supervisor.cron
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=""
HOME=/
*/5 * * * * root ${datadog_name}
